# profile.d/broot.sh
#!/usr/bin/env bash

eval "$(broot --print-shell-function bash)"
