# config/systemd/user/gnome-keyring.service
[Unit]
Description=Keyring
[Service]
ExecStart=/usr/local/bin/export-keys
KillUserProcesses=yes
[Install]
WantedBy=default.target
