# config/systemd/user/gpg-agent.service.d/override.conf
[Service]
Environment="GPG_KEY_ID=4F0BDB86C4684DDD"
Environment="PASSPHRASE=00000000"
